# deployment/attendance.service
[Unit]
Description=Attendance System Gunicorn Daemon
After=network.target mysql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/attendance
Environment="DJANGO_SETTINGS_MODULE=attendance_project.settings_production"
ExecStart=/var/www/attendance/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 attendance_project.wsgi:application
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
